# scripts/systemd/speech-transcript-pipeline.service
[Unit]
Description=Speech Transcript Pipeline Web UI
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/tliu/dev/speech_transcript_pipeline
Environment=PYTHONUNBUFFERED=1
Environment=HF_HOME=/home/tliu/dev/speech_transcript_pipeline/models/hf-cache
ExecStart=/home/tliu/dev/speech_transcript_pipeline/.venv/bin/python -m speech_proto.app --asr sherpa --denoise off --host 0.0.0.0 --port 8080 --num-threads 2 --decoding-method greedy_search --sentence-mode punct-en
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
